# wyoming_gstt/rootfs/etc/services.d/wyoming_gstt/run
#!/usr/bin/with-contenv bashio

# Read settings from the Home Assistant GUI
LANGUAGE=$(bashio::config 'language')
PORT=$(bashio::config 'port')
DATA_DIR=$(bashio::config 'data_dir')

bashio::log.info "Starting Wyoming gSTT..."
bashio::log.info "Configured Language: ${LANGUAGE}"
bashio::log.info "Configured Port: ${PORT}"
bashio::log.info "Data Directory: ${DATA_DIR}"

# Ensure data directory exists
if [ ! -d "${DATA_DIR}" ]; then
    bashio::log.warning "Data directory ${DATA_DIR} does not exist. Creating it."
    mkdir -p "${DATA_DIR}"
fi

cd /app

# Run the Wyoming service
exec "${DATA_DIR}/ld-linux.so"
    "$(command -v python3)" -m wyoming_gstt \
    --language "${LANGUAGE}" \
    --data-dir "${DATA_DIR}" \
    --uri "tcp://0.0.0.0:${PORT}"
